# docker/python/data/supervisor/flask_app.conf
[program:translate_api]
directory = /var/www/translate_api
command = sh -c 'export PYTHONPATH=/var/www && gunicorn app:app -b localhost:5000 --name translate_api --log-level=warn --log-file=-'
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stderr_logfile = /var/log/supervisor/%(program_name)s.log
autostart=true
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/translate_api.log
numprocs=1